# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/authentik-sync.path
[Unit]
Description=Monitor Authentik certificates

[Path]
PathChanged=/etc/letsencrypt/archive/auth.metropolis.nexus
Unit=authentik-sync.service

[Install]
WantedBy=multi-user.target