# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/forgejo-sync.path
[Unit]
Description=Monitor Forgejo certificates

[Path]
PathChanged=/etc/letsencrypt/archive/forgejo.citadel.chi01.metropolis.nexus
Unit=forgejo-sync.service

[Install]
WantedBy=multi-user.target