# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/jitsi-sync.path
[Unit]
Description=Monitor Jitsi certificates

[Path]
PathChanged=/etc/letsencrypt/archive/jitsi.citadel.chi01.metropolis.nexus
Unit=jitsi-sync.service

[Install]
WantedBy=multi-user.target