# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/nextcloud-sync.path
[Unit]
Description=Monitor Nextcloud certificates

[Path]
PathChanged=/etc/letsencrypt/archive/nextcloud.citadel.chi01.metropolis.nexus
Unit=nextcloud-sync.service

[Install]
WantedBy=multi-user.target