# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/nextcloud-sync.service
[Unit]
Description=Synchronize Nextcloud certificates
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
        
ExecStart=/usr/sbin/rsync -aL --chown 101:101 /etc/letsencrypt/live/nextcloud.citadel.chi01.metropolis.nexus/fullchain.pem root@192.168.2.112:/srv/nextcloud/certs
ExecStart=/usr/sbin/rsync -aL --chown 101:101 /etc/letsencrypt/live/nextcloud.citadel.chi01.metropolis.nexus/privkey.pem root@192.168.2.112:/srv/nextcloud/certs

CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateIPC=true
PrivateTmp=true
ProcSubset=pid
ProtectClock=true
ProtectControlGroups=true
ProtectHome=read-only
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
RestrictAddressFamilies=
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RuntimeDirectoryMode=700
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@obsolete