# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/opencloud-sync.path
[Unit]
Description=Monitor OpenCloud certificates

[Path]
PathChanged=/etc/letsencrypt/archive/opencloud.citadel.chi01.metropolis.nexus
Unit=opencloud-sync.service

[Install]
WantedBy=multi-user.target