# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/pelican-sync.path
[Unit]
Description=Monitor Pelican certificates

[Path]
PathChanged=/etc/letsencrypt/archive/pelican.citadel.chi01.metropolis.nexus
Unit=pelican-sync.service

[Install]
WantedBy=multi-user.target