# etc/systemd/system/nginx.citadel.chi01.metropolis.nexus/vaultwarden-sync.path
[Unit]
Description=Monitor Vaultwarden certificates

[Path]
PathChanged=/etc/letsencrypt/archive/vault.citadel.chi01.metropolis.nexus
Unit=vaultwarden-sync.service

[Install]
WantedBy=multi-user.target