# etc/systemd/system/observatory.nyc01.metropolis.nexus/coturn-sync.path
[Unit]
Description=Monitor Coturn certificates

[Path]
PathChanged=/etc/letsencrypt/archive/turn.metropolis.nexus
Unit=coturn-sync.service

[Install]
WantedBy=multi-user.target