# etc/systemd/system/observatory.nyc01.metropolis.nexus/eturnal-sync.path
[Unit]
Description=Monitor Eturnal certificates

[Path]
PathChanged=/etc/letsencrypt/archive/turn.metropolis.nexus
Unit=eturnal-sync.service

[Install]
WantedBy=multi-user.target