# etc/systemd/system/observatory.nyc01.metropolis.nexus/livekit-sync.path
[Unit]
Description=Monitor LiveKit certificates

[Path]
PathChanged=/etc/letsencrypt/archive/livekit.metropolis.nexus
Unit=livekit-sync.service

[Install]
WantedBy=multi-user.target